# docker/crontab
SHELL=/bin/bash
PATH=/home/seluser/venv/bin:/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
TZ=Asia/Jerusalem
PYTHONUNBUFFERED=1

# Daily sync at 02:00 AM
0 2 * * * root cd /app && python3 bin/sync.py --unlimited --no-approval --verbose --safety-mode regular >> /app/logs/daily_sync.log 2>&1 && echo "$(date -Iseconds)" > /app/healthcheck/last_success.txt

# Logrotate daily at 03:00 AM
0 3 * * * root /usr/sbin/logrotate /etc/logrotate.d/gov2db >> /app/logs/logrotate.log 2>&1
