# docker/crontab
SHELL=/bin/bash
PATH=/home/seluser/venv/bin:/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
TZ=Asia/Jerusalem
PYTHONUNBUFFERED=1

# Randomized sync - runs every 12 hours, wrapper ensures 21-34 hour intervals
# The randomized_sync.sh script:
# 1. Checks if 21+ hours passed since last sync
# 2. Adds 0-13 hours random delay
# 3. Uses --no-headless to bypass Cloudflare
0 2,14 * * * root /app/docker/randomized_sync.sh >> /app/logs/daily_sync.log 2>&1

# Logrotate daily at 03:00 AM
0 3 * * * root /usr/sbin/logrotate /etc/logrotate.d/gov2db >> /app/logs/logrotate.log 2>&1
